# initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup-by-partlabel@sysimages.service.d/80-internal-ordering.conf
[Unit]
After=systemd-cryptsetup-by-partlabel@lukshdrs.service
