# initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup-by-partlabel@sysimages.service.d/90-enroll-fido.conf
[Unit]
After=systemd-cryptenroll-convert-empty-to-fido2@dev-disk-by\x2dpartlabel-%i.service
Wants=systemd-cryptenroll-convert-empty-to-fido2@dev-disk-by\x2dpartlabel-%i.service
