# initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup-detached-by-partuuid@.service.d/50-cryptsetup-options.conf
[Service]
Environment="CRYPTSETUP_OPTIONS=luks,fido2-device=auto,tpm2-device=auto,header=%t/spork/luks/%I.luks"
