# initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup-detached-by-partuuid@.service.d/90-initrd-format-disk.conf
[Unit]
After=initrd-format-disk@%i.target
