# initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup-detached-by-partuuid@.service.d/90-mount-lukshdrs.conf
[Unit]
Wants=run-spork-luks.mount
After=run-spork-luks.mount
ConditionPathExists=%t/spork/luks/%I.luks
