# initrd/mkosi.extra/usr/lib/systemd/system/luks-format-detached-by-partuuid@.service.d/bind-to-tpm2.conf
[Service]
Environment="LUKSFORMAT_TPM2_DEVICE=auto"
Environment="LUKSFORMAT_TPM2_PCRS=0+7+14"
