# initrd/mkosi.extra/usr/lib/systemd/system/luks-format-detached-by-partuuid@.service.d/pvcreate.conf
[Service]
# Run pvcreate on all detached header disks after luksFormat.
Environment="LUKSFORMAT_MKFS=pvcreate"
