# initrd/mkosi.extra/usr/lib/systemd/system/luks-format-by-partlabel@lukshdrs.service.d/luksformat-mkfs-ext4.conf
[Service]
Environment="LUKSFORMAT_MKFS=mkfs.ext4 -L '%I'"
