# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup-detached-by-partuuid@.service.d/50-only-detached-headers.conf
[Service]
# Check that the device does not have an inline LUKS header
# This is an extra precaution to avoid accidentally overwriting a device with a file system
ExecCondition=sh -c "! cryptsetup isLuks '/dev/disk/by-%J/%I'"
