# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup-detached-by-partuuid@.service.d/90-mount-lukshdrs.conf
[Unit]
RequiresMountsFor=%t/spork/luks/
ConditionPathExists=%t/spork/luks/%I.luks
