# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/initrd-init-disks-pre.target
[Unit]
Description=Begin initialize encrypted storage
DefaultDependencies=no
After=systemd-pcrphase-initrd.service
After=initrd-auth.target
Wants=initrd-init-disks.target
Before=initrd-init-disks.target
Before=cryptsetup.target
Before=shutdown.target initrd-switch-root.target
Conflicts=shutdown.target initrd-switch-root.target
OnSuccess=initrd-init-disks.target
OnSuccessJobMode=flush
