# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/lukshdrs-cryptsetup@.path
[Unit]
Description=LUKS detached header for device %I
DefaultDependencies=no
Before=paths.target
After=cryptsetup-pre.target
Before=cryptsetup.target
Before=shutdown.target initrd-switch-root.target
Conflicts=shutdown.target initrd-switch-root.target
Before=installer.target
Conflicts=installer.target

[Path]
PathExists=%t/spork/lukshdrs/%I.luks
