# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/run-spork-lukshdrs.automount
[Unit]
Description=Automount for %P
After=blockdev@dev-mapper-lukshdrs.target
BindsTo=dev-mapper-lukshdrs.device
After=dev-mapper-lukshrds.device
