# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/run-spork-sysimages.automount
[Unit]
Description=sysimages
Wants=blockdev@dev-mapper-sysimages.target
After=blockdev@dev-mapper-sysimages.target
BindsTo=dev-mapper-sysimages.device
