# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/wait-for-yubikey.target
[Unit]
Description=Wait for hardware security token
DefaultDependencies=no
Requires=wait-for-yubikey.service
After=wait-for-yubikey.service
Before=cryptsetup-pre.target
Before=shutdown.target initrd-switch-root.target
Conflicts=shutdown.target initrd-switch-root.target
JobTimeoutSec=infinity
JobRunningTimeoutSec=infinity
