# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup@sysimages.service.d/90-enroll-fido.conf
[Unit]
After=firstboot-cryptenroll-fido2@dev-disk-by\x2dpartlabel-sysimages.service
Wants=firstboot-cryptenroll-fido2@dev-disk-by\x2dpartlabel-sysimages.service
