# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup@sysimages.service.d/90-initrd-only.conf
[Unit]
Before=initrd-switch-root.target
Conflicts=initrd-switch-root.target
