# initrd/mkosi.extra/usr/lib/systemd/system/fido2-token.target
[Unit]
Description=Check for security token
DefaultDependencies=no
After=initrd-auth.target
Requires=dev-security-fido.device
After=dev-security-fido.device
Before=cryptsetup-pre.target
Before=shutdown.target initrd-switch-root.target
Conflicts=shutdown.target initrd-switch-root.target
JobTimeoutSec=30
