# initrd/mkosi.extra/usr/lib/systemd/system/run-spork-sysimages.automount
[Unit]
Description=Automount %f
After=blockdev@dev-mapper-luks\x2d%j.target
BindsTo=dev-mapper-luks\x2d%j.device
After=dev-mapper-luks\x2d%j.device
