# initrd/mkosi.extra/usr/lib/systemd/system/run-spork-sysimages.automount
[Unit]
Description=Automount %f
After=blockdev@dev-mapper-%j.target
BindsTo=dev-mapper-%j.device
After=dev-mapper-%j.device
