# initrd/mkosi.extra/usr/lib/systemd/system/run-spork-luks.mount
[Unit]
Description=Mount %f
After=blockdev@dev-mapper-lukshdrs.target
BindsTo=dev-mapper-lukshdrs.device
After=dev-mapper-lukshdrs.device

[Mount]
What=/dev/mapper/lukshdrs
DirectoryMode=0700
LazyUnmount=yes
