# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/luks-format-by-partlabel@lukshdrs.service.d/luksformat-params.conf
[Service]
# Something other than AES to increase the entropy in our security risk.
Environment="LUKSFORMAT_CIPHER=twofish-xts-plain64"
Environment="LUKSFORMAT_CIPHER_KEYSIZE=512"
# Use dm-integrity for extra protection against outside modifications
Environment="LUKSFORMAT_INTEGRITY=hmac-sha256"
Environment="LUKSFORMAT_MKFS=mkfs.xfs -m rmapbt=1 -L '%I'"
